# lottery-hub.service.example
[Unit]
Description=Lottery Data-Science Report
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/u/pro/3004_lot
ExecStart=/home/u/pro/3004_lot/.venv/bin/python /home/u/pro/3004_lot/serve.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
